# task1/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
priority=2

[program:apache2]
command=apachectl -DFOREGROUND
priority=1
