# task2/supervisord2.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
priority=2

[program:python3]
command=python3 ./server2.py
priority=1